# lumos.conf
# Copy to /etc/init/lumos.conf
# run initctl reload-configuration

description "lumos"

start on runlevel [2345]
stop on runlevel [!2345]

exec /usr/local/bin/lumos.py
